# tests/data/autostart-folders/1/autostart/vim.desktop
[Desktop Entry]
Version=1.0
Name=Vim
Exec=/sbin/vim
Terminal=true
Type=Application
